# scripts/com.aebel.sync-claude-code.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aebel.sync-claude-code</string>

    <key>ProgramArguments</key>
    <array>
        <string>REPO_DIR/launchctl/sync_cc_projects_folder.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>REPO_DIR/logs/sync_cc_projects_folder.log</string>

    <key>StandardErrorPath</key>
    <string>REPO_DIR/logs/sync_cc_projects_folder.err</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
